# MAC_ROOT/Library/LaunchDaemons/userprocesses.aria2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>aria2c RPC service</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/aria2c</string>
		<string>--enable-rpc</string>
		<string>--rpc-listen-all</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/opt/aria2c-downloads</string>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/aria2c-rpc-service.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/aria2c-rpc-service-error.log</string>
</dict>
</plist>
